# files/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
user=root

[program:gohugo]
directory=/home/devopsclub
command=hugo server
user=root
